# rust-solana-scanlink/rust-solana-scanlink-main/rust-solana--gemini/solana-copy-trader/deploy/systemd/rust-scanlink.service
[Unit]
Description=Rust Scanlink Trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rust-scanlink
EnvironmentFile=-/home/ubuntu/rust-scanlink/.env
ExecStart=/home/ubuntu/rust-scanlink/rust-scanlink-trader
Restart=always
RestartSec=3
LimitNOFILE=65535
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
